# deployment/systemd/invoice-beat.service
[Unit]
Description=InvoiceAI Celery Beat Scheduler
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User=invoiceai
Group=invoiceai
WorkingDirectory=/var/www/invoice-app/backend
ExecStart=/bin/bash -c "source /var/www/invoice-app/backend/venv/bin/activate && exec celery -A celery_app beat --loglevel=info"

Restart=always
RestartSec=5

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/invoice-app/storage /tmp
PrivateTmp=true

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
